# 2.3/etc/profile.d/aemanagement.sh
#!/bin/sh
export ANDUTTEYEMANAGEMENT_REPOSITORY="/var/jail/aemanagement"
